# os_base/sys_services/aegis-llm-runtime.service
[Unit]
Description=Aegis LLM Runtime (llama.cpp Server)
After=network.target aegis-agent-core.service
Wants=aegis-agent-core.service

[Service]
Type=simple
User=aegis_agent
Group=aegis
WorkingDirectory=/opt/hispanshield/models
ExecStart=/opt/hispanshield/bin/llama-server -m /opt/hispanshield/models/aegis-core-1.5b.gguf -c 4096 --host 127.0.0.1 --port 8080 --threads 4 --api-key "$LLM_API_KEY"
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
CapabilityBoundingSet=
SystemCallFilter=~@mount @raw-io @debug @privileged
AppArmorProfile=/opt.hispanshield.bin.llama-server
AuditMode=yes
IPAddressAllow=127.0.0.1
IPAddressDeny=any

[Install]
WantedBy=multi-user.target
