# os_base/sys_services/aegis-agent-core.service
[Unit]
Description=Aegis Sentinel Engine Orchestrator (Rust)
After=aegis-llm-runtime.service
Wants=aegis-llm-runtime.service

[Service]
Type=simple
# CWE-269 FIX: Do NOT run as root - contradicts Zero-Trust
User=aegis_admin
Group=aegis_admin
WorkingDirectory=/opt/hispanshield/core/rust/aegis-sentinel
ExecStart=/opt/hispanshield/core/rust/target/release/aegis-sentinel
# Additional hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
RestrictSUIDGID=true
ProtectSystem=full
PrivateTmp=true
ReadOnlyPaths=/etc
ReadWritePaths=/var/log/hispanshield /opt/hispanshield/core/data
AmbientCapabilities=
NoNewPrivileges=true
SystemCallFilter=~@mount @raw-io @debug @privileged
AppArmorProfile=/opt.hispanshield.core.rust.aegis-sentinel
AuditMode=yes

[Install]
WantedBy=multi-user.target
