# os_base/sys_services/aegis-agent-core.service
[Unit]
Description=Aegis Sentinel Engine Orchestrator (Rust)
After=aegis-llm-runtime.service
Wants=aegis-llm-runtime.service

[Service]
Type=simple
User=aegis_admin
WorkingDirectory=/opt/hispanshield/core/rust/aegis-sentinel
ExecStart=/opt/hispanshield/core/rust/target/release/aegis-sentinel
ProtectSystem=full
PrivateTmp=true
ReadOnlyPaths=/etc
ReadWritePaths=/var/log/hispanshield /opt/hispanshield/core/data
AmbientCapabilities=
NoNewPrivileges=true
SystemCallFilter=~@mount @raw-io @debug @privileged
AppArmorProfile=/opt.hispanshield.core.rust.aegis-sentinel
AuditMode=yes

[Install]
WantedBy=multi-user.target
