# os_base/sys_services/aegis-agent-core.service
[Unit]
Description=HispanShield OS LLmSecurity Sentinel Engine (Agent Core)
After=aegis-llm-runtime.service aegis-observer.service
Requires=aegis-llm-runtime.service

[Service]
Type=simple
User=root
# El core del agente corre protegido para poder orquestar tareas críticas, PERO tira privilegios antes de interactuar con el LLM
WorkingDirectory=/opt/HispanShield OS LLmSecurity/core/sentinel_engine
ExecStart=/usr/bin/python3 main.py

# Aislamientos (Hardening)
ProtectSystem=full
# Sólo tiene permisos de escritura sobre su base auditora
ReadWritePaths=/var/log/HispanShield OS LLmSecurity /opt/HispanShield OS LLmSecurity/core/data
PrivateTmp=true
NoNewPrivileges=false

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
